# mounted-vagrant.conf
description	"/vagrant"
start on vagrant-mounted

task

script
    /vagrant/setid.sh
end script
